# automation/seia-agent/seia-agent.service
[Unit]
Description=SEIA Agent - Mise à jour automatique du site Synapsea
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /root/SteFlux-Projet-Synapsea/automation/seia-agent/launch_seia.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target

